# crates/agent/templates/tests/tenant_dhcp-relay.conf
# Auto-generated by NICo!
[program: isc-dhcp-relay-default]
command = /usr/sbin/dhcrelay --nl -d  -i p0_if  -i p1_if   -i vlan196  -i vlan185  -U lo 10.217.5.197 10.217.5.197 -r test
autostart = true
autorestart = unexpected
startsecs = 3
startretries = 3
exitcodes = 0
stopsignal = TERM
stopwaitsecs = 3
