# crates/agent/templates/default-nico-dhcp-server.conf
# Auto-generated by NICo!
[program: nico-dhcp-server-default]
command = /var/support/nico-dhcp/bin/nico-dhcp-server {{range .Interfaces}} --interfaces {{.}} {{ end }} --dhcp-config /var/support/nico-dhcp/conf/dhcp.yaml --host-config /var/support/nico-dhcp/conf/host.yaml
autostart = {{ .Autostart }}
autorestart = {{ .Autostart }}
startsecs = 3
startretries = 3
exitcodes = 0
stopsignal = TERM
stopwaitsecs = 3
stdout_logfile = /var/support/nico-dhcp/logs/nico-dhcp-server.log
redirect_stderr = true
